# launchd/ai.agento.lw-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.agento.lw-watchdog</string>

  <!-- Uses launchd-launcher.sh which sources shell profile for all env vars. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@LAUNCHER_SCRIPT@</string>
    <string>@WATCHDOG_SCRIPT@</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>300</integer>

  <key>WorkingDirectory</key>
  <string>@REPO_ROOT@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>@PATH@</string>
  </dict>

  <key>StandardOutPath</key>
  <string>@LOG_FILE@</string>

  <key>StandardErrorPath</key>
  <string>@LOG_FILE@</string>
</dict>
</plist>
