# launchd/ai.agento.health.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  ai.agento.health — unified AO health check and remediation.

  Replaces the previous multi-plist architecture:
    - ai.agento.lifecycle-all  (start-all.sh launcher)
    - ai.agento.lw-watchdog   (watchdog for lifecycle service)
    - com.agentorchestrator.lw-watchdog (duplicate watchdog)
    - com.agentorchestrator.start-all   (start-all at load)
    - com.agentorchestrator.orch-watchdog

  Runs every 5 min. Ensures lifecycle-workers are running for all
  configured projects, starts missing workers, kills orphans, and
  re-bootstraps its own launchd service if deregistered.

  INSTALLATION (two options):

  Option A — via central setup (recommended):
    cd agent-orchestrator && bash scripts/setup.sh
    # or incrementally:
    bash scripts/setup-launchd.sh health

  Option B — manual:
    # Generate from template by replacing @VAR@ placeholders, then:
    ln -sf "$(pwd)/launchd/ai.agento.health.plist" ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/ai.agento.health.plist

  2026-06-25: log path moved from ~/.openclaw/logs/ to ~/.hermes/logs/ via @LOG_FILE@ replacement.

  The plist is generated by scripts/setup-launchd.sh which substitutes:
    @HOME@, @REPO_ROOT@, @LAUNCHER_SCRIPT@, @HEALTH_SCRIPT@, @LOG_FILE@, @PATH@

  To trigger a test run:
    launchctl kickstart -kp gui/$(id -u)/ai.agento.health
-->
<dict>
  <key>Label</key>
  <string>ai.agento.health</string>

  <!-- Uses launchd-launcher.sh which sources shell profile for all env vars. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@LAUNCHER_SCRIPT@</string>
    <string>@HEALTH_SCRIPT@</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- ao-health.sh forks lifecycle-worker children via nohup then exits.
       AbandonProcessGroup prevents launchd from reaping child workers. -->
  <key>AbandonProcessGroup</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>@REPO_ROOT@</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>@PATH@</string>
    <!-- Go binaries (gh) use Security.framework for TLS cert verification in
         interactive sessions but fail with x509:OSStatus -26276 in launchd agents.
         SSL_CERT_FILE tells Go to use the PEM file instead. -->
    <!-- Non-secret defaults for MiniMax and Claude binary path.
         Secrets (API keys) come from shell profile via launchd-launcher.sh. -->
    <key>MINIMAX_ANTHROPIC_BASE_URL</key>
    <string>https://api.minimax.io/anthropic</string>
    <key>MINIMAX_MODEL</key>
    <string>MiniMax-M2.7</string>
    <key>CLAUDE_BINARY</key>
    <string>@CLAUDE_BINARY@</string>
    <key>SSL_CERT_FILE</key>
    <string>/etc/ssl/cert.pem</string>
    <!-- Source-tree CLI for fork-specific commands (ao skeptic verify, etc.)
         that are absent from the published npm binary. Workers inherit this
         env var; skeptic-reviewer.ts uses it to resolve the ao binary. -->
    <key>AO_CLI_PATH</key>
    <string>@REPO_ROOT@/packages/cli/dist/index.js</string>
    <!-- Belt-and-suspenders: even if a future ao-health.sh change forgets to
         pass --no-open when calling `ao start`, this env var forces
         shouldOpenBrowser() to return false. PR-level guard for the recurring
         "localhost:3000 keeps reopening" complaint. See memory:
         feedback_2026-06-22_ao_health_no_open_regression. -->
    <key>AO_NO_OPEN_BROWSER</key>
    <string>1</string>
  </dict>

  <key>StandardOutPath</key>
  <string>@LOG_FILE@</string>

  <key>StandardErrorPath</key>
  <string>@LOG_FILE@</string>

  <key>ThrottleInterval</key>
  <integer>120</integer>
</dict>
</plist>
