# launchd/ai.hermes.launchd-drift-audit.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  ai.hermes.launchd-drift-audit — nightly structural backstop for plist drift

  Detects LaunchAgents at exit 127 (script-missing / wrong container name)
  and posts a single Slack alert to $HERMES_OPS_SLACK_CHANNEL. Cadence: nightly
  at 02:00 local.

  INSTALLATION (per ~/.claude/skills/launchd-plist-template/SKILL.md):
    bash scripts/setup-launchd.sh launchd-drift-audit
  or (manual; mirrors scripts/setup-launchd.sh launchd-drift-audit):
    # Substitute ALL placeholders (HOME, REPO_ROOT, channel, tokens).
    # setup-launchd.sh does this with escape_sed; the manual sed below uses
    # empty defaults per the PR #615 umbrella pattern (empty > wrong).
    REPO_ROOT="$(cd "$(dirname "$0")" && pwd)"
    sed -e "s|@HOME@|$HOME|g" \
        -e "s|@HERMES_OPS_SLACK_CHANNEL@|${HERMES_OPS_SLACK_CHANNEL:-}|g" \
        -e "s|@SLACK_BOT_TOKEN@|${OPENCLAW_STAGING_SLACK_BOT_TOKEN:-${SLACK_BOT_TOKEN:-}}|g" \
        -e "s|@OPENCLAW_STAGING_SLACK_BOT_TOKEN@|${OPENCLAW_STAGING_SLACK_BOT_TOKEN:-${SLACK_BOT_TOKEN:-}}|g" \
        -e "s|@SLACK_USER_TOKEN@|${SLACK_USER_TOKEN:-}|g" \
        launchd/ai.hermes.launchd-drift-audit.plist.template \
        > ~/Library/LaunchAgents/ai.hermes.launchd-drift-audit.plist
    launchctl bootstrap gui/$(id -u) \
      ~/Library/LaunchAgents/ai.hermes.launchd-drift-audit.plist

  Reference: jleechanorg/agent-orchestrator #709
             feedback_2026-06-18_launchd_plist_drift
-->
<dict>
  <key>Label</key>
  <string>ai.hermes.launchd-drift-audit</string>

  <key>ProgramArguments</key>
  <array>
    <string>@HOME@/.hermes/scripts/launchd-env-wrapper.sh</string>
    <string>@HOME@/.hermes/scripts/audit-launchd-drift.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>@HOME@/.hermes/logs/launchd-drift-audit.out.log</string>

  <key>StandardErrorPath</key>
  <string>@HOME@/.hermes/logs/launchd-drift-audit.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>@HOME@</string>
    <key>PATH</key>
    <string>@HOME@/.bun/bin:@HOME@/.local/bin:@HOME@/bin:@HOME@/Library/pnpm:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!-- Source of truth per PR #615 umbrella pattern. launchd-env-wrapper.sh
         reads this and merges with .bashrc defaults; the script falls back to
         it if HERMES_OPS_SLACK_CHANNEL is empty. NEVER hardcode a channel ID
         or token here — render-time substitution only. -->
    <key>HERMES_OPS_SLACK_CHANNEL</key>
    <string>@HERMES_OPS_SLACK_CHANNEL@</string>
    <!-- Slack tokens for alert posting. setup-launchd.sh substitutes
         @SLACK_BOT_TOKEN@ etc. from the rendered config. Empty default is
         the umbrella pattern (PR #615): the audit script resolves tokens
         from env (bashrc / plist) with empty-fallback, never a redacted
         sentinel. -->
    <key>OPENCLAW_STAGING_SLACK_BOT_TOKEN</key>
    <string>@OPENCLAW_STAGING_SLACK_BOT_TOKEN@</string>
    <key>SLACK_BOT_TOKEN</key>
    <string>@SLACK_BOT_TOKEN@</string>
    <key>SLACK_USER_TOKEN</key>
    <string>@SLACK_USER_TOKEN@</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>ThrottleInterval</key>
  <integer>60</integer>
</dict>
</plist>
